# semaphore-ui/rootfs/etc/services.d/semaphore-ui/run
#!/usr/bin/with-contenv sh

/start_semaphore_ui.sh
